# luci-app-shadowsocks-libev/root/etc/uci-defaults/40_luci-shadowsocks-libev
#!/bin/sh

mkdir -p /etc/shadowsocks-libev
/etc/init.d/rpcd reload
chmod 777 /etc/init.d/shadowsocks-libev
exit 0
